# webui/pantheon-webui.service
[Unit]
Description=Pantheon Web UI
Documentation=https://github.com/Duskript/Pantheon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 %h/pantheon/webui/server.py
WorkingDirectory=%h/pantheon/webui
Restart=always
RestartSec=5

# ── Network ────────────────────────────────────────────────────────────────
# Bind to all interfaces (0.0.0.0) to access from other devices on the LAN
# or via Tailscale/WireGuard. Loopback (127.0.0.1) is the safe default when
# you only need local access.
Environment=HERMES_WEBUI_HOST=0.0.0.0
Environment=HERMES_WEBUI_PORT=8787

# ── State directory (sessions, settings, workspaces) ──────────────────────
# Defaults to ~/.hermes/webui — only set this to redirect elsewhere.
# Environment=HERMES_WEBUI_STATE_DIR=%h/.hermes/webui

[Install]
WantedBy=default.target
